# root/etc/systemd/system/openvpn@.service.d/host-to-net.conf
[Service]
ExecStartPre=/etc/openvpn/openvpn-startup
ExecStopPost=/etc/openvpn/openvpn-shutdown
